# {{cookiecutter.app_name}}/supervisord_programs/gunicorn.conf
[program:gunicorn]
directory=/app
command=gunicorn
    {{cookiecutter.app_name}}.app:create_app()
    -b :5000
    -w %(ENV_GUNICORN_WORKERS)s
    -k gevent
    --max-requests=5000
    --max-requests-jitter=500
    --log-level=%(ENV_LOG_LEVEL)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
